# systemd/user/pacman-dump.timer
[Unit]
Description=Regularily dump a list of installed packages

[Timer]
OnBootSec=0m
OnUnitActiveSec=2h

[Install]
WantedBy=timers.target
