# systemd/user/pacman-dump.timer
[Unit]
Description=Regularily dump a list of installed packages

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
